# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ── Nginx (frontend + reverse proxy) ──────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Node BFF (Fastify on port 3000) ──────────────────────
[program:app]
command=node /app/dist/server.js
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="%(ENV_NODE_ENV)s"
priority=20

; ── Python Worker ─────────────────────────────────────────
[program:worker]
command=python -m src.main
directory=/worker
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/worker"
priority=30

; ── tusd (Resumable Upload Server on port 8080) ──────────
[program:tusd]
command=tusd
    -s3-bucket=%(ENV_TUSD_S3_BUCKET)s
    -s3-endpoint=%(ENV_TUSD_S3_ENDPOINT)s
    -s3-object-prefix=
    -hooks-http=http://127.0.0.1:3000/api/webhooks/tus
    -hooks-enabled-events=pre-create,post-finish
    -max-size=10737418240
    -behind-proxy
    -port=8080
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=AWS_ACCESS_KEY_ID="%(ENV_MINIO_ACCESS_KEY)s",AWS_SECRET_ACCESS_KEY="%(ENV_MINIO_SECRET_KEY)s",AWS_REGION="us-east-1"
priority=25
